# airootfs/etc/systemd/system/custom.service
[Service]
ExecStart=
ExecStart=/usr/bin/start

[Install]
WantedBy=multi-user.target
